# apps/youtube_downloader/youtube-downloader.service
[Unit]
Description=YouTube Downloader Service (External API)
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/youtube-downloader
Environment="YOUTUBE_DOWNLOADER_API_KEY=540cb0e8711ecf1e772b6270f97ea81836bf4a7f2eb5186b9dd34bb70a612e55"
Environment="PORT=8080"
Environment="HOST=0.0.0.0"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/opt/youtube-downloader/venv/bin/python /opt/youtube-downloader/main.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/youtube-downloader/output.log
StandardError=append:/var/log/youtube-downloader/error.log

[Install]
WantedBy=multi-user.target
